# apps/ui/sys-services/in-contact-ui.service
[Unit]
Description="In Contact UI"

[Service]
ExecStart=/usr/bin/node --experimental-modules --es-module-specifier-resolution=node app
WorkingDirectory=/home/ubuntu/inContact/apps/ui
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=in-contact-ui
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
